# bin/mute.crontab
0  8  * * * ${HOME}/bin/unmute.sh >> ${HOME}/bin/mute.log
30 21 * * * ${HOME}/bin/mute.sh >> ${HOME}/bin/mute.log

