# pause-sonos.crontab
0 22 * * * ${HOME}/bin/pause-sonos.sh >> ${HOME}/bin/pause-sonos.log
